# src/backend/project-s.service
[Unit]
Description=Project S Backend API
After=network.target mysql.service

[Service]
Type=simple
User=sean
WorkingDirectory=/home/sean/project-s/src/backend
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
